# ffserver/rootfs/etc/services.d/httpd/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting HTTPD..."

# Run
exec /bin/busybox-extras httpd -f -v -p 8080 -h /tmp